# .claude/skills/scheduler/templates/recurring-macos.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.schedule.{ID}</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "{WORKING_DIR}" &amp;&amp; {CLAUDE_PATH} -p "{PROMPT}" {FLAGS}</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:{HOME}/.local/bin</string>
        <key>HOME</key>
        <string>{HOME}</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>{HOUR}</integer>
        <key>Minute</key><integer>{MINUTE}</integer>
        <!-- Add Weekday for specific days: <key>Weekday</key><integer>3</integer> -->
    </dict>
    <key>StandardOutPath</key>
    <string>{HOME}/.claude/logs/{ID}.log</string>
    <key>StandardErrorPath</key>
    <string>{HOME}/.claude/logs/{ID}.error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
